# system/systemd/open-notebook-worker.service
[Unit]
Description=Open Notebook Worker
PartOf=open-notebook.target
After=open-notebook-surreal.service open-notebook-api.service
Requires=open-notebook-surreal.service

[Service]
Type=simple
ExecStart=/usr/bin/env bash __PROJECT_DIR__/system/open-notebook-systemd-run.sh worker
SyslogIdentifier=open-notebook-worker
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=3
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=open-notebook.target
